# deploy/stock-agent.service
[Unit]
Description=Stock Data Service (FastAPI/uvicorn)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# 以下路径由 remote_deploy.sh 用实际值渲染后安装到 /etc/systemd/system/
User=__RUN_USER__
WorkingDirectory=__APP_DIR__/service
EnvironmentFile=__APP_DIR__/.env
Environment=CACHE_DIR=__APP_DIR__/cache
Environment=DATA_DIR=__APP_DIR__/data
ExecStart=__VENV__/bin/uvicorn app:app --host 0.0.0.0 --port __PORT__
Restart=on-failure
RestartSec=3
# 优雅停止，给 uvicorn 时间释放端口
KillSignal=SIGINT
TimeoutStopSec=15
# ---- 内存护栏（由 remote_deploy.sh 按总内存渲染，可用 MEMORY_MAX/MEMORY_HIGH 覆盖）----
# MemoryHigh：软上限，超过后内核对本服务施压回收；MemoryMax：硬上限，超过即在本 cgroup
# 内 OOM，只杀本服务、随后 Restart=on-failure 自动拉起，避免拖垮整机与 SSH。
MemoryHigh=__MEMORY_HIGH__
MemoryMax=__MEMORY_MAX__
OOMPolicy=kill

[Install]
WantedBy=multi-user.target
